# etc-init.d/graphics-driver
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# 
# Copyright (c) 2019 Triffid Hunter
# 
# /etc/init.d/graphics-driver
# automatically load appropriate nvidia userspace library which matches kernel version
#
# for best results:
#     list /lib/modules in CONFIG_PROTECT so emerge doesn't remove modules from other kernels
#     emerge nvidia-drivers && service graphics-driver restore when updating your drivers and kernel
#     make xdm depend on this

depend() {
	need localmount
	before xdm
}

start() {
	NVPKG="$(strings /lib/modules/$(uname -r)/video/nvidia.ko | grep x11-drivers/nvidia-drivers | head -n1 | grep -o 'x11-drivers.*' | cut -d/ -f1-2)"
	if ! [ -d "/var/db/pkg/$NVPKG" ]
	then
		echo "Want $NVPKG, have $(echo /var/db/pkg/x11-drivers/nvidia-driver* | cut -d/ -f5-6)"
		if emerge -u1KO "=$NVPKG" && [ -d /sys/module/nvidia ] && ! nvidia-smi
		then
			rmmod nvidia_drm nvidia_modeset nvidia && modprobe nvidia
		fi
	fi
}

stop() {
	emerge -u1KO nvidia-drivers
}
